# files/etc/systemd/system/kubelet.service.d/timeout.conf
# Give kubelet plenty of time to terminate gracefully, communicating NodeShutdown statues for all pods.
# This should be higher than shutdownGracePeriod.
[Service]
TimeoutStopSec=180s
