# files-kiosk/etc/systemd/system/dpms@off.timer
[Unit]
Description=Turn screen off

[Timer]
OnCalendar=*-*-* 00:30:00

[Install]
WantedBy=timers.target
