# files-kiosk/etc/systemd/system/dpms@on.timer
[Unit]
Description=Turn screen on

[Timer]
OnCalendar=*-*-* 09:00:00

[Install]
WantedBy=timers.target
